# chloe-blog.service
[Unit]
Description=Chloe's Blog Auto-Update
Documentation=https://github.com/stellarie/chloe-blog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/stella/chloe-blog/update_blog.sh
WorkingDirectory=/home/stella/chloe-blog

# Log output
StandardOutput=journal
StandardError=journal

# Run as user
User=stella

[Install]
WantedBy=default.target
