# chloe-blog.timer
[Unit]
Description=Chloe's Blog Auto-Update Timer
Documentation=https://github.com/stellarie/chloe-blog
Requires=chloe-blog.service

[Timer]
# Run every hour at :00
OnCalendar=hourly
OnBootSec=5min
Persistent=true

# Randomize slightly to avoid thundering herd
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
